# basic-mesh/files/etc/hotplug.d/iface/30-basic-mesh
#!/bin/sh
#
# /etc/hotplug.d/iface/30-basic-mesh
#
# Applies 802.11s mesh_param settings from UCI when a mesh interface comes up.
#
# There are two categories of configured mesh interfaces:
#
# "Claimed" — the ifname appears in a wifi-iface section in /etc/config/wireless
#   with mode=mesh. For these, we know which network they belong to, so we only
#   apply params when the matching network's ifup fires.
#
# "Unclaimed" — the ifname is configured in basic-mesh but does not appear in
#   any wifi-iface ifname (e.g. auto-generated names like wlan0-mesh discovered
#   via iw). For these, we apply params on every ifup event since we cannot tie
#   them to a specific network.

logger -t basic-mesh-hotplug "triggered: ACTION=${ACTION} INTERFACE=${INTERFACE} DEVICE=${DEVICE}"

[ "${ACTION}" = "ifup" ] || exit 0

# Build a list of all wifi-iface ifnames with mode=mesh, and apply params
# for any whose network matches the hotplug INTERFACE.
# Returns the set of "claimed" ifnames via the CLAIMED variable.
CLAIMED=""

for section in $(uci -q show wireless | sed -n 's/^wireless\.\([^.]*\)=wifi-iface$/\1/p'); do
	wifi_mode=$(uci -q get "wireless.${section}.mode")
	[ "$wifi_mode" = "mesh" ] || continue

	ifname=$(uci -q get "wireless.${section}.ifname")
	if [ -z "$ifname" ]; then
		logger -t basic-mesh-hotplug "wireless.${section} has mode=mesh but no ifname set, cannot claim"
		continue
	fi

	CLAIMED="${CLAIMED} ${ifname}"

	wifi_net=$(uci -q get "wireless.${section}.network")
	: "${wifi_net:=lan}"

	[ "$wifi_net" = "${INTERFACE}" ] || continue

	bm_type=$(uci -q get "basic-mesh.${ifname}")
	if [ "$bm_type" != "mesh_params" ]; then
		logger -t basic-mesh-hotplug "claimed ${ifname} (wireless.${section}, network=${wifi_net}) — no basic-mesh config, skipping"
		continue
	fi

	logger -t basic-mesh-hotplug "claimed ${ifname} (wireless.${section}, network=${wifi_net}) — applying"
	basic-mesh-apply "$ifname"
done

# Now handle unclaimed basic-mesh sections: apply on every ifup.
for bm_iface in $(uci -q show basic-mesh | sed -n 's/^basic-mesh\.\([^.]*\)=mesh_params$/\1/p'); do
	# Skip if this ifname is claimed by a wifi-iface section
	case " ${CLAIMED} " in
		*" ${bm_iface} "*) continue ;;
	esac

	logger -t basic-mesh-hotplug "unclaimed ${bm_iface} — applying on ifup ${INTERFACE}"
	basic-mesh-apply "$bm_iface"
done

exit 0
