# basic-mesh/files/etc/hotplug.d/iface/30-basic-mesh
#!/bin/sh
#
# /etc/hotplug.d/iface/30-basic-mesh
#
# Applies 802.11s mesh_param settings from UCI when a mesh interface comes up.
#
# The hotplug INTERFACE is the network-level name (e.g. 'mesh') and DEVICE is
# typically the bridge (e.g. 'br-mesh'). Neither is the actual wireless
# interface (e.g. 'mesh0'). To find it, we look up wifi-iface sections in
# /etc/config/wireless where mode=mesh and network matches INTERFACE, then
# use the ifname from that section.

logger -t basic-mesh-hotplug "triggered: ACTION=${ACTION} INTERFACE=${INTERFACE} DEVICE=${DEVICE}"

[ "${ACTION}" = "ifup" ] || exit 0

# Find mesh wifi-iface sections whose 'network' matches the hotplug INTERFACE,
# resolve their ifname, and apply basic-mesh params if configured.
apply_mesh_ifaces() {
	local net="$1"
	[ -z "$net" ] && return

	local found=0
	local section ifname wifi_net wifi_mode

	for section in $(uci -q show wireless | sed -n 's/^wireless\.\([^.]*\)=wifi-iface$/\1/p'); do
		wifi_mode=$(uci -q get "wireless.${section}.mode")
		[ "$wifi_mode" = "mesh" ] || continue

		wifi_net=$(uci -q get "wireless.${section}.network")
		# network defaults to 'lan' if unset
		: "${wifi_net:=lan}"

		[ "$wifi_net" = "$net" ] || continue

		ifname=$(uci -q get "wireless.${section}.ifname")
		if [ -z "$ifname" ]; then
			logger -t basic-mesh-hotplug "wireless.${section} has mode=mesh, network=${net} but no ifname set, skipping"
			continue
		fi

		logger -t basic-mesh-hotplug "wireless.${section}: network=${net} -> ifname=${ifname}"
		found=1

		local bm_section
		bm_section=$(uci -q get "basic-mesh.${ifname}")
		if [ -z "$bm_section" ]; then
			logger -t basic-mesh-hotplug "no basic-mesh config for '${ifname}', skipping"
			continue
		fi

		if [ "$bm_section" != "mesh_params" ]; then
			logger -t basic-mesh-hotplug "basic-mesh.${ifname} is type '${bm_section}', not mesh_params, skipping"
			continue
		fi

		logger -t basic-mesh-hotplug "applying mesh params for '${ifname}'"
		basic-mesh-apply "$ifname"
	done

	[ "$found" -eq 0 ] && logger -t basic-mesh-hotplug "no wireless mesh interface found for network '${net}'"
}

apply_mesh_ifaces "${INTERFACE}"

exit 0
