# basic-mesh/files/etc/uci-defaults/99-basic-mesh
#!/bin/sh
# Ensure the UCI config file exists (may have been removed manually)
touch /etc/config/basic-mesh
exit 0
